# docker-file/supervisord.conf
# cat supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/sbin/httpd -k start

[program:redis]
command = /usr/bin/redis-server
